# src/etc/NetworkManager/dispatcher.d/99-autohotspot
#!/bin/bash
INTERFACE=$1
STATUS=$2

# Only run for wlan0
if [ "$INTERFACE" != "wlan0" ]; then
    exit 0
fi

if [ "$STATUS" = "up" ]; then
    # If up, check if we are on Home_WiFi or Hotspot
    CURRENT_CON=$(nmcli -t -f NAME connection show --active | grep "Home_WiFi")
    if [ -n "$CURRENT_CON" ]; then
        /usr/bin/logger "Autohotspot: Connected to Home WiFi."
    fi
elif [ "$STATUS" = "down" ] || [ "$STATUS" = "disconnected" ]; then
    /usr/bin/logger "Autohotspot: WiFi Down. Checking for Home Network..."
    # Check if Home SSID is visible
    if nmcli device wifi list | grep -q "Your_Home_SSID"; then
        /usr/bin/logger "Autohotspot: Home network found. Connecting..."
        nmcli connection up "Home_WiFi"
    else
        /usr/bin/logger "Autohotspot: Home network NOT found. Starting Hotspot."
        nmcli connection up "Car_Hotspot"
    fi
fi
