# src/etc/systemd/system/bluetooth.service.d/override.conf
[Service]
StateDirectory=
ReadWritePaths=/mnt/data/bluetooth
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPost=/usr/local/bin/bt-startup.sh
